# srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon=true
user=root
[program:influxdb]
command=/usr/sbin/influxd -config /etc/influxdb.conf
autorestart=true
[program:telegraf]
command=./telegraf/usr/bin/telegraf
autorestart=true